# misc/upstart/eve_static.conf
author "Brian Hechinger <wonko@4amlunch.net>"
start on runlevel [23]

script
  cd /home/ubuntu/es
  ./eve_static
end script
